# DEBIAN/postinst
#!/bin/bash
set -e

# ============================================================================
# Bad IPs Package Post-Installation Script
# Implements privilege separation by running service as non-root user
# ============================================================================

# Postinst parameters:
# $1 = action (configure, abort-upgrade, abort-remove, abort-deconfigure)
# $2 = old-version (if upgrading) or version being installed (if configure)
ACTION="$1"
OLD_VERSION="$2"

# ============================================================================
# UPGRADE DETECTION: Handle upgrades from version 3.4.7 or older
# ============================================================================
UPGRADING_FROM_OLD_VERSION=0

if [ "$ACTION" = "configure" ] && [ -n "$OLD_VERSION" ]; then
    echo "Detected upgrade from version $OLD_VERSION"

    # Check if upgrading from version < 3.5.0 (versions without privilege separation)
    # Version comparison: 3.4.7 < 3.5.0
    if dpkg --compare-versions "$OLD_VERSION" lt "3.5.0"; then
        UPGRADING_FROM_OLD_VERSION=1
        echo ""
        echo "╔════════════════════════════════════════════════════════════════════════╗"
        echo "║                   IMPORTANT: PRIVILEGE SEPARATION UPGRADE              ║"
        echo "╚════════════════════════════════════════════════════════════════════════╝"
        echo ""
        echo "This version introduces privilege separation for enhanced security."
        echo "The bad_ips service will now run as a non-root user instead of root."
        echo ""
        echo "Changes being applied:"
        echo "  • Service will run as: badips user (non-root)"
        echo "  • Limited sudo access: Only nftables inet badips table operations"
        echo "  • Config files: Secured with proper permissions"
        echo "  • The service will be restarted to apply these changes"
        echo ""

        # Check if config file doesn't exist - need to create it for upgrade
        if [ ! -f /etc/default/bad_ips ]; then
            echo "Creating service user configuration with default username: badips"
            echo ""

            cat > /etc/default/bad_ips <<'EOF'
# Bad IPs service user configuration
# This file is automatically generated during upgrade
# The bad_ips service will run as this user for security (privilege separation)
#
# To change the username after installation:
#   1. Stop the service: systemctl stop bad_ips
#   2. Edit this file and change BADIPS_USER and BADIPS_GROUP
#   3. Run: /usr/local/share/bad_ips/reconfigure-user.sh
#   4. Start the service: systemctl start bad_ips

BADIPS_USER=badips
BADIPS_GROUP=badips
EOF
            chmod 644 /etc/default/bad_ips

            echo "→ Service user configuration created: /etc/default/bad_ips"
            echo "→ Default username: badips"
            echo ""
            echo "If you want to use a different username, you can change it later by:"
            echo "  1. Editing /etc/default/bad_ips"
            echo "  2. Running: /usr/local/share/bad_ips/reconfigure-user.sh"
            echo ""
        fi

        # Stop the service before making ownership changes
        if systemctl is-active --quiet bad_ips.service; then
            echo "Stopping bad_ips service for privilege separation upgrade..."
            systemctl stop bad_ips.service || true
            NEED_RESTART=1
        fi
    fi
fi

# Load service user configuration (created by installer or use default)
if [ -f /etc/default/bad_ips ]; then
    . /etc/default/bad_ips
fi
BADIPS_USER=${BADIPS_USER:-badips}
BADIPS_GROUP=${BADIPS_GROUP:-badips}

echo "Configuring Bad IPs service to run as user: $BADIPS_USER"

# ============================================================================
# 1. Create service user if it doesn't exist (fallback)
# ============================================================================
if ! getent passwd "$BADIPS_USER" > /dev/null; then
    echo "Creating system user: $BADIPS_USER"
    if adduser --system --group --no-create-home \
        --home /nonexistent --shell /usr/sbin/nologin \
        "$BADIPS_USER" > /dev/null 2>&1; then
        echo "✓ System user '$BADIPS_USER' created"
    else
        echo "ERROR: Failed to create user '$BADIPS_USER'"
        exit 1
    fi
else
    echo "✓ User '$BADIPS_USER' already exists"
fi

# ============================================================================
# 2. Add user to supplementary groups for log access
# ============================================================================
# Add to systemd-journal group (for journalctl access)
if getent group systemd-journal > /dev/null 2>&1; then
    usermod -aG systemd-journal "$BADIPS_USER" 2>/dev/null || true
fi

# Add to adm group (for /var/log access on Debian/Ubuntu)
if getent group adm > /dev/null 2>&1; then
    usermod -aG adm "$BADIPS_USER" 2>/dev/null || true
fi

# ============================================================================
# 3. Create and configure directories with proper ownership
# ============================================================================
# Create database directory
if [ ! -d /var/lib/bad_ips ]; then
    mkdir -p /var/lib/bad_ips
fi
chown -R "$BADIPS_USER:$BADIPS_GROUP" /var/lib/bad_ips
chmod 755 /var/lib/bad_ips

# Create log directory and file
if [ ! -d /var/log/bad_ips ]; then
    mkdir -p /var/log/bad_ips
fi
chown -R "$BADIPS_USER:$BADIPS_GROUP" /var/log/bad_ips
chmod 755 /var/log/bad_ips

if [ ! -f /var/log/bad_ips/bad_ips.log ]; then
    touch /var/log/bad_ips/bad_ips.log
fi
chown "$BADIPS_USER:$BADIPS_GROUP" /var/log/bad_ips/bad_ips.log
chmod 640 /var/log/bad_ips/bad_ips.log

# Create cache directory for public blocklists
if [ ! -d /var/cache/badips ]; then
    mkdir -p /var/cache/badips
fi
chown -R "$BADIPS_USER:$BADIPS_GROUP" /var/cache/badips
chmod 755 /var/cache/badips

# ============================================================================
# 4. Secure configuration files (contain database passwords)
# ============================================================================
# Main config file - allow service user to read (root-owned)
if [ -f /usr/local/etc/badips.conf ]; then
    chmod 640 /usr/local/etc/badips.conf
    chown root:"$BADIPS_GROUP" /usr/local/etc/badips.conf
fi

# Config directory and all .conf files
if [ -d /usr/local/etc/badips.d ]; then
    chmod 750 /usr/local/etc/badips.d
    chown root:"$BADIPS_GROUP" /usr/local/etc/badips.d

    # Secure all config files in the directory
    find /usr/local/etc/badips.d -type f -name "*.conf" -exec chmod 640 {} \; 2>/dev/null || true
    find /usr/local/etc/badips.d -type f -name "*.conf" -exec chown root:"$BADIPS_GROUP" {} \; 2>/dev/null || true
fi

echo "✓ File ownership and permissions configured"

# ============================================================================
# 4.5. Fix application file permissions (world-readable for service user)
# ============================================================================
echo "Ensuring correct permissions on application files..."

# Perl modules must be world-readable for service user to load them
if [ -d /usr/local/lib/site_perl/BadIPs ]; then
    chmod -R a+rX /usr/local/lib/site_perl/BadIPs
fi

# log4perl configuration must be readable by service user
if [ -f /usr/local/etc/badips/log4perl.conf ]; then
    chmod 644 /usr/local/etc/badips/log4perl.conf
fi

# Ensure main script is executable
if [ -f /usr/local/sbin/bad_ips ]; then
    chmod 755 /usr/local/sbin/bad_ips
fi

echo "✓ Application file permissions configured"

# ============================================================================
# 5. Create sudoers configuration for nftables access
# ============================================================================
echo "Creating sudoers configuration for $BADIPS_USER..."

cat > /etc/sudoers.d/bad_ips <<SUDOEOF
# Bad IPs sudoers configuration
# Allows $BADIPS_USER user to manage the inet badips nftables table
# Created by bad-ips package installation
# DO NOT EDIT - Changes will be overwritten on package upgrade

# Command aliases for allowed nft operations
Cmnd_Alias NFT_BADIPS_ADD = /usr/sbin/nft add element inet badips badipv4 *
Cmnd_Alias NFT_BADIPS_ADD_V6 = /usr/sbin/nft add element inet badips badipv6 *
Cmnd_Alias NFT_BADIPS_FLUSH = /usr/sbin/nft flush set inet badips never_block*, \\
                              /usr/sbin/nft flush set inet badips always_block*
Cmnd_Alias NFT_BADIPS_LIST = /usr/sbin/nft -j list ruleset, \\
                             /usr/sbin/nft list ruleset

# Allow $BADIPS_USER user to run these commands without password
$BADIPS_USER ALL=(root) NOPASSWD: NFT_BADIPS_ADD, NFT_BADIPS_ADD_V6, NFT_BADIPS_FLUSH, NFT_BADIPS_LIST
SUDOEOF

# Set correct permissions on sudoers file
chmod 0440 /etc/sudoers.d/bad_ips

# Validate sudoers syntax
if ! visudo -c -f /etc/sudoers.d/bad_ips >/dev/null 2>&1; then
    echo "ERROR: Invalid sudoers syntax in /etc/sudoers.d/bad_ips"
    echo "Removing invalid file to prevent system issues"
    rm -f /etc/sudoers.d/bad_ips
    exit 1
fi

echo "✓ Sudoers configuration created and validated"

# ============================================================================
# 6. Create systemd drop-in for user configuration
# ============================================================================
echo "Creating systemd service configuration..."

mkdir -p /etc/systemd/system/bad_ips.service.d

cat > /etc/systemd/system/bad_ips.service.d/user.conf <<SYSDEOF
# Bad IPs service user configuration
# Created by bad-ips package installation
# This configures the service to run as a non-root user for security

[Service]
# Run as non-root user with limited privileges
User=$BADIPS_USER
Group=$BADIPS_GROUP

# Supplementary groups for log/journal access
SupplementaryGroups=systemd-journal adm
SYSDEOF

echo "✓ Systemd drop-in configuration created"

# Setup nftables (if service is active)
# The nftables configuration is managed by /etc/nftables.d/99-badips.nft
# This file is created by the installer script and loaded automatically
if systemctl is-active --quiet nftables.service; then
    echo "Reloading nftables to apply Bad IPs configuration..."
    systemctl reload nftables.service || true
else
    echo "Note: nftables.service is not active - Bad IPs rules not loaded"
    echo "      Enable nftables to activate Bad IPs blocking"
fi

# Run database migration if database exists
if [ -f /var/lib/bad_ips/bad_ips.sql ]; then
    echo "Checking database schema version..."
    if [ -f /usr/local/share/bad_ips/migrate_database_v2.pl ]; then
        perl /usr/local/share/bad_ips/migrate_database_v2.pl --db-path /var/lib/bad_ips/bad_ips.sql
        if [ $? -eq 0 ]; then
            echo "Database schema up to date"
        else
            echo "WARNING: Database migration failed. Check logs."
        fi
    fi
fi

# Reload systemd daemon
systemctl daemon-reload

# Check if config file exists
if [ ! -f /usr/local/etc/badips.conf ]; then
    echo ""
    echo "==================================================================="
    echo "IMPORTANT: Configuration Required"
    echo "==================================================================="
    echo ""
    echo "Bad IPs has been installed but requires configuration."
    echo ""
    echo "Please copy the template to /usr/local/etc/badips.conf:"
    echo ""
    echo "  sudo cp /usr/local/etc/badips.conf.template /usr/local/etc/badips.conf"
    echo ""
    echo "CRITICAL: Edit the config and set 'never_block_cidrs' with your trusted networks!"
    echo ""
    echo "After configuration:"
    echo "  sudo systemctl enable bad_ips.service"
    echo "  sudo systemctl start bad_ips.service"
    echo ""
    echo "==================================================================="
    echo ""
else
    echo "Configuration file already exists at /usr/local/etc/badips.conf"
    echo "To upgrade: review changes in the template files if needed"

    # ============================================================================
    # UPGRADE: Add PublicBlocklistPlugins section if not present (v3.5.0+)
    # ============================================================================
    if ! grep -q "^\[PublicBlocklistPlugins:" /usr/local/etc/badips.conf 2>/dev/null; then
        echo ""
        echo "Adding new PublicBlocklistPlugins configuration section..."

        # Backup the config file first
        cp /usr/local/etc/badips.conf /usr/local/etc/badips.conf.bak-v3.5.0
        echo "  • Backup created: /usr/local/etc/badips.conf.bak-v3.5.0"

        # Append the PublicBlocklistPlugins section
        cat >> /usr/local/etc/badips.conf <<'PLUGINEOF'

# Public Blocklist Plugin Configurations (added in v3.5.0)
# (optional) urls - Comma separated list of URLs to fetch the blocklist from
# (optional) fetch_interval - How often (in seconds) to fetch the blocklist
# (optional) use_cache - 1 to use caching, 0 to not use caching
# (optional) cache_path - Directory path to store cached files (if use_cache is 1)
# (required*) active - 1 to enable this blocklist plugin, 0 to disable it
# Each plugin MUST have
# * a unique name after PublicBlocklistPlugins:
# * at least active = 1 to be used; all other parameters are really at the use of the plugin
[PublicBlocklistPlugins:Spamhaus]
urls = https://www.spamhaus.org/drop/drop.txt, https://www.spamhaus.org/drop/edrop.txt
fetch_interval = 3600
use_cache = 1
cache_path = /var/cache/badips/
active = 1

[PublicBlocklistPlugins:Feodotracker]
urls = https://feodotracker.abuse.ch/downloads/ipblocklist.txt
fetch_interval = 7200
use_cache = 1
cache_path = /var/cache/badips/
active = 0

[PublicBlocklistPlugins:Blocklist_de]
urls = https://lists.blocklist.de/lists/all.txt
fetch_interval = 10800
use_cache = 1
cache_path = /var/cache/badips/
active = 0

[PublicBlocklistPlugins:DNSBL_Info]
urls = https://www.dnsbl.info/dnsbl-list.php?file=dnslist-ipv4.txt
fetch_interval = 14400
use_cache = 1
cache_path = /var/cache/badips/
active = 0

[PublicBlocklistPlugins:Malwaredomainlist]
urls = http://www.malwaredomainlist.com/hostslist/ip.txt
fetch_interval = 21600
use_cache = 1
cache_path = /var/cache/badips/
active = 0
PLUGINEOF

        # Restore proper permissions
        chmod 640 /usr/local/etc/badips.conf
        chown root:"$BADIPS_GROUP" /usr/local/etc/badips.conf

        echo "  • PublicBlocklistPlugins section added (Spamhaus enabled by default)"
        echo "  • To disable: Set active = 0 for any plugin in badips.conf"
        echo ""
    fi

    # On upgrade, restore service state
    if [ -f /var/lib/bad_ips/.upgrade_was_enabled ]; then
        echo "Re-enabling bad_ips.service (was enabled before upgrade)"
        systemctl enable bad_ips.service
        rm -f /var/lib/bad_ips/.upgrade_was_enabled
    fi

    if [ -f /var/lib/bad_ips/.upgrade_was_active ]; then
        echo "Restarting bad_ips.service (was active before upgrade)"
        systemctl start bad_ips.service
        rm -f /var/lib/bad_ips/.upgrade_was_active
    fi
fi

# ============================================================================
# UPGRADE: Restart service if it was stopped for privilege separation upgrade
# ============================================================================
if [ "$UPGRADING_FROM_OLD_VERSION" = "1" ] && [ "${NEED_RESTART:-0}" = "1" ]; then
    echo ""
    echo "╔════════════════════════════════════════════════════════════════════════╗"
    echo "║              PRIVILEGE SEPARATION UPGRADE COMPLETE                     ║"
    echo "╚════════════════════════════════════════════════════════════════════════╝"
    echo ""
    echo "Security improvements applied:"
    echo "  ✓ Service user '$BADIPS_USER' created and configured"
    echo "  ✓ Sudo access limited to nftables inet badips table only"
    echo "  ✓ Configuration files secured with proper permissions"
    echo "  ✓ SystemD configured to run as non-root user"
    echo ""
    echo "Starting bad_ips service with new security configuration..."

    if systemctl start bad_ips.service; then
        echo "✓ Service started successfully as user: $BADIPS_USER"
        echo ""
        echo "You can verify the service is running with:"
        echo "  systemctl status bad_ips.service"
        echo "  ps aux | grep bad_ips"
        echo ""
    else
        echo "✗ Failed to start service"
        echo ""
        echo "Please check the logs:"
        echo "  journalctl -u bad_ips.service -n 50"
        echo ""
        echo "Common issues:"
        echo "  • Check sudo configuration: sudo -u $BADIPS_USER sudo nft list ruleset"
        echo "  • Check file permissions: ls -la /var/log/bad_ips/ /var/lib/bad_ips/"
        echo "  • Check user groups: groups $BADIPS_USER"
        echo ""
    fi

    echo "For more information about privilege separation, see:"
    echo "  /usr/share/doc/bad-ips/README.md"
    echo ""
fi

exit 0
