# DEBIAN/postinst
#!/bin/bash
set -e

# Create database directory
if [ ! -d /var/lib/bad_ips ]; then
    mkdir -p /var/lib/bad_ips
    chmod 755 /var/lib/bad_ips
fi

# Create log directory and file
if [ ! -d /var/log/bad_ips ]; then
    mkdir -p /var/log/bad_ips
    chmod 755 /var/log/bad_ips
fi

if [ ! -f /var/log/bad_ips/bad_ips.log ]; then
    touch /var/log/bad_ips/bad_ips.log
    chmod 640 /var/log/bad_ips/bad_ips.log
fi

# Setup nftables (if service is active)
# The nftables configuration is managed by /etc/nftables.d/99-badips.nft
# This file is created by the installer script and loaded automatically
if systemctl is-active --quiet nftables.service; then
    echo "Reloading nftables to apply Bad IPs configuration..."
    systemctl reload nftables.service || true
else
    echo "Note: nftables.service is not active - Bad IPs rules not loaded"
    echo "      Enable nftables to activate Bad IPs blocking"
fi

# Run database migration if database exists
if [ -f /var/lib/bad_ips/bad_ips.sql ]; then
    echo "Checking database schema version..."
    if [ -f /usr/local/share/bad_ips/migrate_database_v2.pl ]; then
        perl /usr/local/share/bad_ips/migrate_database_v2.pl --db-path /var/lib/bad_ips/bad_ips.sql
        if [ $? -eq 0 ]; then
            echo "Database schema up to date"
        else
            echo "WARNING: Database migration failed. Check logs."
        fi
    fi
fi

# Reload systemd daemon
systemctl daemon-reload

# Check if config file exists
if [ ! -f /usr/local/etc/badips.conf ]; then
    echo ""
    echo "==================================================================="
    echo "IMPORTANT: Configuration Required"
    echo "==================================================================="
    echo ""
    echo "Bad IPs has been installed but requires configuration."
    echo ""
    echo "Please copy the appropriate template to /usr/local/etc/badips.conf:"
    echo ""
    echo "For hunter mode (most servers):"
    echo "  sudo cp /usr/local/etc/badips.conf.hunter-template /usr/local/etc/badips.conf"
    echo ""
    echo "For gatherer mode (administrator server):"
    echo "  sudo cp /usr/local/etc/badips.conf.gatherer-template /usr/local/etc/badips.conf"
    echo ""
    echo "CRITICAL: Edit the config and set 'never_block_cidrs' with your trusted networks!"
    echo ""
    echo "After configuration:"
    echo "  sudo systemctl enable bad_ips.service"
    echo "  sudo systemctl start bad_ips.service"
    echo ""
    echo "==================================================================="
    echo ""
else
    echo "Configuration file already exists at /usr/local/etc/badips.conf"
    echo "To upgrade: review changes in the template files if needed"

    # On upgrade, restore service state
    if [ -f /var/lib/bad_ips/.upgrade_was_enabled ]; then
        echo "Re-enabling bad_ips.service (was enabled before upgrade)"
        systemctl enable bad_ips.service
        rm -f /var/lib/bad_ips/.upgrade_was_enabled
    fi

    if [ -f /var/lib/bad_ips/.upgrade_was_active ]; then
        echo "Restarting bad_ips.service (was active before upgrade)"
        systemctl start bad_ips.service
        rm -f /var/lib/bad_ips/.upgrade_was_active
    fi
fi

exit 0
